# claude-automation-pipeline/scripts/com.id8labs.ava.memory-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.id8labs.ava.memory-sync</string>

    <key>Comment</key>
    <string>Ava Memory Sync - Cross-platform memory sync between Telegram and Parallax every 30 minutes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/your-username/.hydra/daemons/ava-memory-sync.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/your-username/.local/bin</string>
        <key>HOME</key>
        <string>/Users/your-username</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/your-username/Library/Logs/claude-automation/ava-memory-sync/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/your-username/Library/Logs/claude-automation/ava-memory-sync/launchd-stderr.log</string>

    <key>TimeOut</key>
    <integer>120</integer>

    <key>RunAtLoad</key>
    <false/>
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
